# docker-nginx-node-2016-tw-election/my_nodejs_2016_tw_election.conf
description "twitter container"
author "Alan Tai"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a prj2016twelection_node_1
end script
